# crontab
MAILTO="meteor_a_key@yahoo.co.jp"

*/5 * * * * echo "`date`"

*/30 * * * * /Users/tatsuya/workspace_ruby/nicoliver

